# deploy/polysiem.service
[Unit]
Description=PolySIEM homelab documentation dashboard
Documentation=https://github.com/Realynx/PolySIEM
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=polysiem
Group=polysiem
WorkingDirectory=/opt/polysiem/run
EnvironmentFile=/opt/polysiem/.env
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0
# Generated/uploaded web certificates survive release-bundle swaps here.
Environment=POLYSIEM_CERT_DIR=/opt/polysiem/data/certs
ExecStart=/usr/bin/node tls-server.js
Restart=always
RestartSec=5

# Hardening
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
